# deploy/systemd/hass-sma-rs.service
[Unit]
Description=hass-sma-rs — SMA Sunny Boy BT daemon
Documentation=https://github.com/MaximV93/hass-sma-rs
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/bin/hass-sma-daemon --config /etc/hass-sma-rs/config.yaml
Restart=on-failure
RestartSec=10
User=hassma
Group=hassma
# BlueZ raw sockets need CAP_NET_RAW + CAP_NET_ADMIN.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
DevicePolicy=closed

[Install]
WantedBy=multi-user.target
